# deployments/systemd/exchange-data-platform-update.timer
[Unit]
Description=Nightly rebuild and restart for Exchange Data Platform

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=exchange-data-platform-update.service

[Install]
WantedBy=timers.target
