# deployments/systemd/exchange-data-platform-backup.service
[Unit]
Description=Exchange Data Platform lake backup
After=exchange-data-platform.service

[Service]
Type=oneshot
WorkingDirectory=/opt/exchange-data-platform
EnvironmentFile=/opt/exchange-data-platform/deployments/docker/.env
ExecStart=/usr/bin/bash /opt/exchange-data-platform/scripts/backup/backup_lake.sh

[Install]
WantedBy=multi-user.target
